# .config/systemd/user/vibmax-server.service
[Unit]
Description=Vibmax Orchestration Server
After=network.target
Wants=opencode.service

[Service]
Type=simple
ExecStart=/home/user/.local/bin/vibmax-server
Restart=on-failure
RestartSec=5
Environment=OPENCODE_URL=http://localhost:4096
Environment=VIBMAX_PORT=4097
Environment=OPENCODE_EXPERIMENTAL_HTTPAPI=true

[Install]
WantedBy=default.target
